# omnibot-v2.5.service
[Unit]
Description=OMNIBOT v2.5 Intelligent Trading System
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=Omnibot
Group=Omnibot
WorkingDirectory=/home/Omnibot/omnibot/src
Environment="PYTHONPATH=/home/Omnibot/omnibot/src"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/Omnibot/omnibot/venv/bin/python /home/Omnibot/omnibot/src/main.py --mode cli
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
